# private_dot_config/systemd/user/home-disk-trend.service
[Unit]
Description=Track home disk growth trend and alert on large growth

[Service]
Type=oneshot
ExecStart=/home/pips/Developer/docker/scripts/home-chatops-host-details.sh trend
