# private_dot_config/systemd/user/home-disk-trend.timer
[Unit]
Description=Track home disk growth trend periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Persistent=true
AccuracySec=5m
Unit=home-disk-trend.service

[Install]
WantedBy=timers.target
