# private_dot_config/systemd/user/home-chatops-host-details.timer
[Unit]
Description=Poll PipsBot ChatOps host-side detail requests

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=10s
Unit=home-chatops-host-details.service

[Install]
WantedBy=timers.target
