# private_dot_config/systemd/user/ps2-vmc-retention.service
[Unit]
Description=Apply safe PS2 VMC backup retention

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/pips/Developer/docker/scripts/ps2-vmc-retention.py --apply
